# examples/current-implementation/systemd/zcs-energy-poller.service
# /etc/systemd/system/zcs-energy-poller.service
[Unit]
Description=ZCS Energy MQTT Poller
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=home-automation
WorkingDirectory=/opt/home-automation/energy-poller
ExecStartPre=/bin/sleep 20
ExecStart=/opt/home-automation/energy-poller/venv/bin/python /opt/home-automation/energy-poller/zcs_energy_poller.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
